# overlay-auto-update.service
[Unit]
Description=OverlayFS Safe Auto Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/flock -n /var/run/overlay-auto-update.lock /usr/local/sbin/overlay-auto-update.sh
RemainAfterExit=no
Restart=on-failure
